# conf/neta-inbox.conf
# /etc/init/neta-inbox.conf
description "neta-inbox"
author  "daisuke-shimada"

start on runlevel [2345]
stop on runlevel [016]

exec /opt/neta-inbox/bin/neta-inbox 2>&1

respawn
